# Display/services.conf
[program:publish_server]
directory=/src/zmq_proxy
command=python3 server.py

[program:zmq_proxy]
directory=/src/zmq_proxy
command=./zmq_proxy